# plow-airbnb-kiosk.service
[Unit]
Description=Plow Airbnb Calendar Kiosk Mode
After=network-online.target graphical.target plow-airbnb-dashboard.service
Wants=network-online.target

[Service]
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/odio/.Xauthority
Type=simple
User=odio
ExecStartPre=/bin/sleep 5
ExecStartPre=/bin/rm -f /home/odio/.config/chromium/Singleton*
ExecStart=/usr/bin/chromium --kiosk --noerrdialogs --disable-infobars --no-first-run --check-for-update-interval=31536000 http://localhost:5174
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
